# kiosk-demo/files/etc/systemd/system/kiosk.service
[Unit]
Description=Edge Kiosk Mode
After=graphical.target

[Service]
User=kiosk
Environment=DISPLAY=:0
ExecStart=/etc/kiosk/kiosk.sh
Restart=always

[Install]
WantedBy=graphical.target

